# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn views:app -b localhost:5000
directory=/deploy/app
